# haana/rootfs/etc/s6-overlay/s6-rc.d/svc-qdrant/run
#!/command/with-contenv bashio
# Qdrant Vector Database

bashio::log.info "Starting Qdrant..."

mkdir -p /data/qdrant

exec /usr/local/bin/qdrant \
    --storage-path /data/qdrant \
    2>&1
