# haana/rootfs/etc/s6-overlay/s6-rc.d/svc-qdrant/run
#!/command/with-contenv bashio
# Qdrant Vector Database

bashio::log.info "Starting Qdrant..."

mkdir -p /data/qdrant

export QDRANT__STORAGE__STORAGE_PATH="/data/qdrant"
export QDRANT__SERVICE__HTTP_PORT="6333"
export QDRANT__TELEMETRY_DISABLED="true"

exec /usr/local/bin/qdrant 2>&1
